# hooks/openclaw/tokenleague-openclaw-collector.timer
[Unit]
Description=Run TokenLeague OpenClaw collector every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=tokenleague-openclaw-collector.service

[Install]
WantedBy=timers.target
