# hooks/openclaw/tokenleague-openclaw-collector.service
[Unit]
Description=TokenLeague OpenClaw collector
After=network-online.target
Wants=network-online.target
ConditionPathExists={{COLLECTOR_PATH}}

[Service]
Type=oneshot
User={{OPENCLAW_USER}}
Environment=HOME={{OPENCLAW_HOME}}
Environment=OPENCLAW_BIN_PATH={{OPENCLAW_BIN_PATH}}
Environment=PATH={{OPENCLAW_BIN_DIR}}:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory={{OPENCLAW_HOME}}
ExecStart={{PYTHON_BIN}} {{COLLECTOR_PATH}}
